# nodeserver.service
[Unit]
Description=REST API service
After=cloud-final.service
Wants=cloud-final.service

[Service]
Type=simple
User=admin
ExecStart=node /home/admin/web2/app.js
EnvironmentFile=/home/admin/web2/.env
StandardOutput=syslog
StandardError=syslog
Restart=on-failure

[Install]
WantedBy=cloud-init.target